# systemd/camilladsp-autoswitch.service
[Unit]
Description=CamillaDSP Autoswitch Daemon
After=network.target camilladsp.service
Requires=camilladsp.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m camilladsp_autoswitch.autoswitch
Restart=always

[Install]
WantedBy=multi-user.target
